# dist/linux/systemd/dbforbix-nojsvc.service
[Unit]
Description=DBforBIX Monitoring Agent
After=syslog.target network.target

[Service]
Type=simple
User=<DBFORBIX_USER>
Group=<DBFORBIX_GROUP>
Environment=JAVA_HOME=<JAVA_HOME_PATH>
WorkingDirectory=/opt/dbforbix
ExecStart=<JAVA_HOME_PATH>/bin/java -jar dbforbix.jar -a start -C /opt/dbforbix
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
